# launcher/launcher.service
[Unit]
Description=Launcher (turn on/off relay)

[Service]
Type=forking
ExecStart=/home/root/launcher/launcher
Environment="HOME=/home/root"
